# mnt-john.mount
[Unit]
Description=automount
[Mount]
What= /dev/mmcblk2p1
Where=/mnt/john/
[Install]
WantedBy=multi-user.target
